# Atlass.Framework.Web/Atlass.Framework.Web.ini
[program:Atlass.Framework.Web]
command=/www/wwwroot/root/easycms/Atlass.Framework.Web
directory=/www/wwwroot/root/easycms/
autorestart=true
startsecs=3
startretries=3
stdout_logfile=/www/server/panel/plugin/supervisor/log/Atlass.Framework.Web.out.log
stderr_logfile=/www/server/panel/plugin/supervisor/log/Atlass.Framework.Web.err.log
stdout_logfile_maxbytes=2MB
stderr_logfile_maxbytes=2MB
user=root
priority=999
numprocs=1
process_name=%(program_name)s_%(process_num)02d